# pkg/utils/templates/tmpws.service.tmpl
[Unit]
After=network.target
[Service]
Type=simple
WorkingDirectory=/opt
ExecStart=/usr/bin/python3 -m http.server {{.WSPort}}
[Install]
WantedBy=default.target
